# aws_s3_scripts/move_daily_tick_data.cron
0 3 * * * source ~/.bash_profile && /home/timir/move_daily_tick_data.sh FXCM

